# backroot8/rootfs-overlay/etc/systemd/system/backroot8-fb-splash.service
[Unit]
Description=Backroot 8 framebuffer boot emblem (before Plymouth)
DefaultDependencies=no
After=local-fs.target systemd-vconsole-setup.service
Before=plymouth-start.service backroot8-desktop.service
Conflicts=getty@tty1.service

[Service]
Type=simple
Environment=TERM=linux
ExecStart=/usr/share/backroot8/show-fb-splash.sh
ExecStop=-/usr/bin/killall -q fbi
StandardInput=tty
StandardOutput=tty
StandardError=null

[Install]
WantedBy=sysinit.target
