# backroot8/rootfs-overlay/etc/systemd/system/backroot8-fb-splash.service
[Unit]
Description=Backroot 8 framebuffer boot emblem (Plymouth fallback)
DefaultDependencies=no
After=systemd-vconsole-setup.service plymouth-start.service
Before=backroot8-desktop.service
Conflicts=getty@tty1.service

[Service]
Type=simple
Environment=TERM=linux
ExecStart=/usr/share/backroot8/show-fb-splash.sh
ExecStop=-/usr/bin/killall -q fbi
TTYPath=/dev/tty1
TTYVTDisallocate=no
StandardInput=null
StandardOutput=null
StandardError=journal

[Install]
WantedBy=multi-user.target
