# backroot8/rootfs-overlay/etc/profile.d/backroot8.sh
export PS1='[\u@Backroot8 \W]\$ '
export XCURSOR_THEME=Win7OS-cursors
export XCURSOR_SIZE=24
